# config/gammu-supervisor.conf
[program:sms_gammu]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/tp-sms-api/artisan gammu:async:handle
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/tp-sms-api/storage/logs/sms.log